# recipes-init/openrc/openrc-scripts/bluetooth.initd
#!/sbin/openrc-run

BTD=/usr/sbin/bluetoothd
PID=/var/run/bluetoothd.pid

depend() {
	before xdm
	need dbus localmount hostname
}

start() {
	ebegin "Enabling Bluetooth"
	start-stop-daemon --start --background --make-pidfile --pidfile $PID --exec $BTD
	eend $?
}

stop() {
	start-stop-daemon --stop --pidfile $PID
	eend $?
}

# vim: noet ft=gentoo-init-d
