# recipes-init/openrc/openrc-scripts/busybox-ntpd.initd
#!/sbin/openrc-run

depend() {
	need net
	use dns logger
}

command_args="-n"
command_background="true"
pidfile="/var/run/busybox-ntpd.pid"
description="busybox ntpd"
command="/usr/sbin/ntpd"

# vim: noet ft=gentoo-init-d
