# recipes-init/openrc/openrc-scripts/connman.initd
#!/sbin/openrc-run

depend() {
	need dbus
	provide net
	after wpa_supplicant
}

command_args="${CONNMAN_OPTS}"
description="connection manager"
command="/usr/sbin/connmand"

# vim: noet ft=gentoo-init-d
