# recipes-appends/oe-core/apmd/apmd.initd
#!/sbin/openrc-run

command="/usr/sbin/apmd"
command_args=${COMMAND_ARGS}
name="APM Daemon"

depend() {
	need localmount
	use logger
}

# vim: noet ft=gentoo-init-d
