# recipes-appends/oe-core/bluez5/bluetooth.initd
#!/sbin/openrc-run

command=/usr/libexec/bluetooth/bluetoothd
command_args=${COMMAND_ARGS}
pidfile=/var/run/bluetoothd.pid
command_background="true"

depend() {
	after coldplug
	need dbus localmount hostname
}

# vim: noet ft=gentoo-init-d
