# recipes-appends/oe-core/at/atd.initd
#!/sbin/openrc-run

command="/usr/sbin/atd"
command_args="-f ${COMMAND_ARGS}"
command_background="true"
pidfile="/var/run/atd.pid"

depend() {
	need clock
	use logger
}

# vim: noet ft=gentoo-init-d
