# recipes-appends/oe-core/busybox/busybox-mdev.initd
#!/sbin/openrc-run

depend() {
	need sysfs dev-mount
	before checkfs fsck
}

command_args="-df ${COMMAND_ARGS}"
command_background="true"
pidfile="/var/run/busybox-mdev.pid"
command="/sbin/mdev"

# vim: noet ft=gentoo-init-d
