# recipes-appends/oe-core/tcf-agent/tcf-agent.initd
#!/sbin/openrc-run

extra_started_commands="reload"

command="/usr/sbin/tcf-agent"
command_args=${COMMAND_ARGS}
command_background="true"
pidfile="/var/run/tcf-agent.pid"
stopsig="USR2"

depend() {
	need localmount
	after net
}

# vim: noet ft=gentoo-init-d
